# systemdServices/dashboard.service
[Unit]
Description=LED Dashboard
After=network.target
Wants=network.target

[Service]
# Run as root for hardware access (LED matrix requires GPIO permissions)
User=root
WorkingDirectory=/opt/led-dashboard

# Start the TypeScript source directly with bun
ExecStart=/usr/bin/bun run src/index.ts

# Auto-restart on crashes
Restart=always
RestartSec=5

# Environment file for optional settings (safe if missing)
EnvironmentFile=-/opt/led-dashboard/.env

# Where logs go
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
