# systemdServices/update@.service
[Unit]
Description=LED Dashboard Update Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/led-dashboard
ExecStart=/opt/led-dashboard/scripts/update-with-service-control.sh %i
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
